# etc/local.d/vfio-bind.start
#!/bin/sh

# Bind GPU's IOMMU group
/usr/local/bin/vfio-bind 0000:01:00.0 0000:01:00.1
